# watchdog/alphadivision-watchdog.service
[Unit]
Description=AlphaDivision Watchdog
Documentation=https://github.com/nickchow0/alphadivision
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu

WorkingDirectory=/opt/alphadivision
ExecStart=/usr/bin/python3 /opt/alphadivision/watchdog/watchdog.py

# Restart automatically on any failure; wait 10 s before retrying so a
# crash loop (e.g. Redis down at boot) doesn't spin at full speed.
Restart=always
RestartSec=10

# Give the watchdog 30 s to shut down cleanly before systemd kills it.
TimeoutStopSec=30

# Log to the systemd journal (view with: journalctl -u alphadivision-watchdog -f)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=alphadivision-watchdog

# Environment — overrides defaults baked into watchdog.py.
# All secrets come from /opt/alphadivision/.env (loaded by the script itself
# via python-dotenv), so we only set the compose dir here.
Environment=COMPOSE_DIR=/opt/alphadivision

[Install]
WantedBy=multi-user.target
